# arbiter3/scripts/config_templates/arbiter-eval.service.jinja
[Unit]
Description=Arbiter Evaluation Loop
After=network.target

[Service]
Type=simple
User={{ user }}
Environment=DJANGO_SETTINGS_MODULE="settings"
WorkingDirectory={{ working_dir }}
ExecStart={{ python }} {{ working_dir }}/arbiter.py evaluate --seconds 30
ExecReload=/bin/kill -s HUP $MAINPID